# config/systemd/subyard-keys-sync.service.in
[Unit]
Description=Synchronize enrolled Subyard encrypted credential ledgers
Documentation=man:systemd.service(5)

[Service]
Type=oneshot
ExecStart=@YARD_BIN@ _keys-auto-sync --if-due

# The service handles encrypted records only. Keep a conservative process sandbox while allowing
# writes to the operator-owned Subyard config/state roots and outbound SSH to enrolled owner hosts.
NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
